# build/osman-autostart.service
[Unit]
Description=Osman starts VM tagged with auto
Requires=virtqemud.service virtnetworkd.service virtqemud-admin.socket
After=virtqemud.service virtnetworkd.service virtqemud-admin.socket

[Service]
Type=oneshot
Environment="HOME=/root"
ExecStart=/bin/sh -c "/usr/bin/osman drop --type=vm --all && /usr/bin/osman start :auto"
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
